# build/logio/supervisor_server.conf
[supervisord]
nodaemon=false
logfile=/dev/stdout
logfile_maxbytes=10MB
loglevel=warn
loglevel=error
user=logio
logfile=/home/logio/supervisord/supervisord.log
pidfile=/home/logio/supervisord/supervisord.pid
childlogdir=/tmp/

[program:logio]
command=log.io-server &
stdout_logfile = /tmp/logio-server.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false