# build/logio/supervisor_harvester.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=10MB
loglevel=warn
loglevel=error
user=logio
logfile=/home/logio/supervisord/supervisord.log
pidfile=/home/logio/supervisord/supervisord.pid
childlogdir=/tmp/

[program:harvester]
command=bash /home/logio/harvester_conf.sh
stdout_logfile = /tmp/harvester.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false