# build/logio/supervisor_harvester.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=warn
user=root
logfile=/home/logio/supervisord.log
pidfile=/var/run/supervisord.pid

loglevel=error

[program:harvester]
command=bash /home/logio/harvester_conf.sh