# build/logio/supervisor_harvester.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=10MB
loglevel=warn
user=logio
logfile=/home/logio/supervisord/supervisord.log
pidfile=/home/logio/supervisord/supervisord.pid
childlogdir=/tmp/

[program:harvester]
command=bash /home/logio/harvester_conf.sh
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr